# server/anni-website.service
# ──────────────────────────────────────────────────────────
#  systemd service for AnniBackend
#  Place at: /etc/systemd/system/anni-website.service
#  Then:
#    sudo systemctl daemon-reload
#    sudo systemctl enable anni-website
#    sudo systemctl start anni-website
#    sudo systemctl status anni-website
#
#  The service runs from /opt/anni/server on the SD card so
#  the backend process itself is always available. The
#  organizer.db inside db/ is a symlink to the storage drive;
#  if that drive is unplugged the backend handles the missing
#  DB gracefully (see server.js degraded-mode handling).
# ──────────────────────────────────────────────────────────

[Unit]
Description=AnniBackend — OAuth API for yumehana.dev
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=akira
WorkingDirectory=/opt/anni/server
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Load secrets from .env file
EnvironmentFile=/opt/anni/server/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anni-website

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
# Needs write access to the local db/ dir (sessions.db + the organizer.db symlink target)
# and transitively to /srv/storage where the organizer.db actually lives.
ReadWritePaths=/opt/anni/server /srv/storage

[Install]
WantedBy=multi-user.target
