# stats/anni-stats.service
[Unit]
Description=anni-stats — System stats API for yumehana.dev
After=network.target

[Service]
Type=simple
User=akira
WorkingDirectory=/srv/storage/AnniWebsite/stats
ExecStart=/usr/bin/python3 stats.py
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=anni-stats

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
